# deploy/systemd/fastapi-daily-incremental-retry.timer
[Unit]
Description=Retry Pending Daily Incremental Sync Every Hour

[Timer]
OnCalendar=*-*-* *:05:00
Persistent=true
Unit=fastapi-daily-incremental-retry.service

[Install]
WantedBy=timers.target
